# intelcs-streaming-agent.service
[Unit]
Description=Intel CS streaming-agent
After=syslog.target network.target
[Service]
Type=forking
PIDFile=/opt/intelcs/production/logs/streaming-agent.pid
User=root
Group=root
WorkingDirectory=/opt/intelcs/production
ExecStart=/opt/intelcs/production/bin/daemon.sh start streaming-agent
ExecStop=/opt/intelcs/production/bin/daemon.sh stop streaming-agent
Restart=always
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNOFILE=64000
LimitNPROC=64000
[Install]
WantedBy=multi-user.target